# templates/default/confd.conf.erb
# Generated by Chef.
description "confd is a lightweight configuration management tool"

start on started etcd
stop on stopping etcd

exec /usr/local/bin/confd <%= @args %>
